# starlamp.service
[Unit]
Description=Star Lamp Service

[Service]
Type=simple
ExecStart=/usr/bin/starlamp

[Install]
WantedBy=multi-user.target